# systemd/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn.socket

SocketUser=wyleu

[Install]
WantedBy=sockets.target
