# systemd/gunicorn.service
[Unit]
Description=Gunicorn daemon for Sandbells
Requires=gunicorn.socket
After=network.target

[Service]
Type=notify
User=wyleu
Group=wyleu

RuntimeDirectory=gunicorn
WorkingDirectory=/home/wyleu/Code/Sandbells/changes

# Optimized for Raspberry Pi 3
ExecStart=/home/wyleu/Code/Sandbells/Bellvirtenv/bin/gunicorn \
    --access-logfile - \
    --workers 1 \
    --threads 4 \
    --worker-class gthread \
    --max-requests 800 \
    --max-requests-jitter 100 \
    --timeout 60 \
    --keep-alive 2 \
    --log-level warning \
    --bind unix:/run/gunicorn.socket \
    changes.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

# Pi 3 friendly
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
